# deploy/systemd/Nika-AI-web.service
[Unit]
Description=Nika-AI Next.js PWA
After=network.target Nika-AI-api.service
Requires=Nika-AI-api.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/Nika-AI/apps/web
EnvironmentFile=/opt/Nika-AI/.env
Environment=ENVIRONMENT=production
Environment=API_PROXY_URL=http://127.0.0.1:8000
Environment=NODE_ENV=production
ExecStart=/usr/bin/pnpm start
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
